# scripts/.config/systemd/user/backup-repos.timer
[Unit]
Description=Daily Development backup timer

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=10min

[Install]
WantedBy=timers.target
